# event-jukebox-audio.service
[Unit]
Description=Event Jukebox Audio Service
After=network.target
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/event-jukebox
ExecStart=/usr/bin/node /home/pi/event-jukebox/audio-integration.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=NODE_ENV=production
Environment=JUKEBOX_STATUS_LOG=0
SyslogIdentifier=jukebox-audio

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
