# wedding-jukebox-audio.service
[Unit]
Description=Wedding Jukebox Audio Service
After=network.target
Requires=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/wedding-jukebox
ExecStart=/usr/bin/env node audio-integration.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=NODE_ENV=production

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
